# etc/systemd/system/mysqld.service.d/systemd-limits.conf
[Service]
LimitNOFILE=65535
LimitMEMLOCK=infinity
